# etc/gap-website.path
[Unit]
Description="Check whether GAP webhook requested update"

[Path]
PathChanged=/home/gap-www/gap-website.trigger
Unit=gap-website.service

[Install]
WantedBy=multi-user.target
