# etc/gap-website.path
[Unit]
Description="Check whether GAP webhook requested update"

[Path]
PathChanged=/srv/www/www-gap-dev/data/gap-website.trigger
Unit=gap-website.service

[Install]
WantedBy=default.target
