# etc/gap-website.service
[Unit]
Description="Update GAP website"

[Service]
User=gap-www
Group=gap-www
Type=oneshot
ExecStart=/home/gap-www/GapWWW/etc/update.sh

[Install]
WantedBy=multi-user.target
